# rootfs/etc/inittab
null::sysinit:/bin/mount -t proc proc /proc 
null::sysinit:/bin/mkdir -p /dev/pts 
null::sysinit:/bin/mount -a 
null::sysinit:/sbin/ifconfig lo 127.0.0.1 up 
::sysinit:/etc/rc 
tty1::respawn:/sbin/getty 38400 tty1 
null::sysinit:/bin/touch /var/log/messages 
::ctrlaltdel:/sbin/reboot 
null::shutdown:/usr/bin/killall klogd 
null::shutdown:/usr/bin/killall syslogd 
null::shutdown:/bin/umount -a -r"
